# rpi/services/rpi-bot.service
[Unit]
Description=RPi Learning Companion
After=network.target bluetooth.target sound.target
Wants=bluetooth.target

[Service]
Type=simple
User=hairobo
Group=hairobo
WorkingDirectory=/home/hairobo/rpi-bot/rpi
ExecStart=/home/hairobo/rpi-bot/rpi/.venv/bin/python -m core.main
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/hairobo/rpi-bot/rpi/data

[Install]
WantedBy=multi-user.target
