# rpi/services/rpi-bot.service
[Unit]
Description=RPi Learning Companion
After=network.target bluetooth.target sound.target
Wants=bluetooth.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/rpi-bot/rpi
ExecStart=/home/pi/rpi-bot/rpi/.venv/bin/python -m core.main
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/pi/rpi-bot/rpi/data

[Install]
WantedBy=multi-user.target
